# sentinel/labs/ansible/roles/sentinel/templates/sentinel.service.j2
[Unit]
Description=Sentinel {{ sentinel_version }}
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ sentinel_service_user }}
Group={{ sentinel_service_group }}
ExecStart=/usr/bin/python3 {{ sentinel_entrypoint }} --config {{ sentinel_config_path }} serve
ExecStartPre=/usr/bin/python3 {{ sentinel_entrypoint }} --config {{ sentinel_config_path }} --check-config
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ sentinel_state_directory }}

[Install]
WantedBy=multi-user.target
