# sentinel/labs/sentinel-app/checkpoints/1.1.0/packaging/sentinel.service
[Unit]
Description=Sentinel Security Platform
Documentation=file:///usr/share/doc/sentinel/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
NotifyAccess=main
User=sentinel
Group=sentinel
ExecStartPre=/usr/libexec/sentinel/sentinel --config /etc/sentinel/sentinel.conf --check-config
ExecStart=/usr/libexec/sentinel/sentinel --config /etc/sentinel/sentinel.conf serve
ExecStartPost=/usr/libexec/sentinel/sentinel --config /etc/sentinel/sentinel.conf --healthcheck
WatchdogSec=10s
TimeoutStopSec=15s
Restart=on-failure
RestartSec=5s
StateDirectory=sentinel
StateDirectoryMode=0750
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
